# build/DEBIAN/postinst
#!/bin/bash

echo "install"

systemctl enable app
systemctl start app
