# build/DEBIAN/postrm
#!/bin/bash

echo "hello world remove"

systemctl disable app
