# build/DEBIAN/postrm
#!/bin/bash

echo "hello world remove"

systemctl stop app
systemctl disable app
systemctl daemon-reload
